# hub/scripts/aws/cron.conf
[program:cron]
command = sudo cron -f -L 15
autostart=true
autorestart=true